# zbserver.conf
[program:device-controller]
user=root
command=/usr/bin/zbserver
autostart = true
autorestart = true

